# packages/network/bluez/system.d/bluetooth-mac.service
[Unit]
Description=Bluetooth MAC configuration

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/setbtmac
StandardInput=tty
TTYPath=/dev/tty1

[Install]
WantedBy=bluetooth.target
